# files/etc/crontabs/www-data
# run website maintenance task
# min   hour    day     month   weekday command
*/1     *       *       *       *       php /var/www/html/public/wp/wp-cron.php
